# hid-io.service
[Service]
ExecStart=/usr/bin/hid-io-core
DynamicUser=yes

