# hid-io.service
[Service]
ExecStart=/usr/bin/hid-io
DynamicUser=yes

